# bin/cron/sys.cron
# system cron. 
#min   hour day month week  command. 
0      1    *   *     *     D:\kenan.zhang\Downloads\ZenTaoPMS_v8.2.5\zentaopms\bin\dailyreminder.bat # daily reminder.           
1      1    *   *     *     D:\kenan.zhang\Downloads\ZenTaoPMS_v8.2.5\zentaopms\bin\backup.bat        # backup database and file. 
1      23   *   *     *     D:\kenan.zhang\Downloads\ZenTaoPMS_v8.2.5\zentaopms\bin\computeburn.bat   # compute burndown chart.   
1-59/2 *    *   *     *     D:\kenan.zhang\Downloads\ZenTaoPMS_v8.2.5\zentaopms\bin\syncsvn.bat       # sync subversion.          
1-59/2 *    *   *     *     D:\kenan.zhang\Downloads\ZenTaoPMS_v8.2.5\zentaopms\bin\syncgit.bat       # sync git.                 
